# services/svd-apache-nc.ini
[supervisord]
nodaemon=true

[program:control]
command=/usr/bin/python3.9 -u /usr/local/bin/control_server.py
directory=/tmp
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=5
priority=5

[program:dispatcher]
command=/usr/local/bin/wait_httpd.sh __default__
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=999
priority=20

[program:renderer]
command=/usr/local/bin/wait_httpd.sh /etc/httpd/conf/httpd-renderer.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=999
priority=20
