# services/svd-apache-nc.ini
[supervisord]
nodaemon=true

[program:renderer]
command=/usr/sbin/httpd -D FOREGROUND -f /etc/httpd/conf/httpd-renderer.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dispatcher]
command=/usr/sbin/httpd -D FOREGROUND
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
